# docker/supervisord.conf
[supervisord]
nodaemon=true  # 启动后不进入守护进程模式

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'  # 启动 Nginx 服务
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log  # 错误日志
stdout_logfile=/var/log/nginx.out.log  # 输出日志

[program:probe-server]
command=/app/probe-server   # 启动 Go 后端服务
autostart=true
autorestart=true
stderr_logfile=/var/log/probe-server.err.log
stdout_logfile=/var/log/probe-server.out.log
